# root/etc/s6-overlay/s6-rc.d/svc-de/run
#!/usr/bin/with-contenv bash

CUSER="${CUSTOM_USER:-chromium}"
HOME="/home/${CUSER}"
cd $HOME
exec s6-setuidgid ${CUSER} \
  /bin/bash /defaults/startwm.sh
